# apps/linkedin-leads/infra/systemd/linkedin-leads-external-ping.timer
[Unit]
Description=Run linkedin-leads-external-ping every 5 minutes
Documentation=https://github.com/Unobtainiumrock/linkedin-leads

[Timer]
OnBootSec=2min
OnUnitActiveSec=5min
Unit=linkedin-leads-external-ping.service
# Catch up if the laptop was asleep when a tick was supposed to fire.
Persistent=true

[Install]
WantedBy=timers.target
